# cron_scripts/crontab
0 3 * * * /home/weather/script 18
0 8 * * * /home/weather/script 0
0 13 * * * /home/weather/script 6
0 18 * * * /home/weather/script 12
